# docker/conf/supervisor.ini
[program:lighttpd]
command         = /usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
redirect_stderr = true
priority        = 600
startretries    = 10

[program:daemon]
command         = python3 /daemon.py
redirect_stderr = true
priority        = 600
startretries    = 10