# docker/conf/supervisor.ini
[program:lighttpd]
command         = /usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
redirect_stderr = true
priority        = 600
startretries    = 10

[program:daemon]
command         = python3 /var/www/localhost/htdocs/daemon.py
redirect_stderr = true
priority        = 600
startretries    = 10

[program:daemon]
command         = aria2c --conf-path=/aria2.conf
redirect_stderr = true
priority        = 600
startretries    = 10

